# docker/whitelotto-php/config/supervisor/whitelotto-cron.conf
[program:cron]
process_name=%(program_name)s
command=/usr/sbin/cron -f
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/php/cron-server.log
