# docker/whitelotto-php/config/cron.d/whitelotto
* * * * * www-data php /var/www/platform/oil r scheduler >> /dev/null 2>&1
